# deploy/Yunshi_supervisor.ini
[program:Yunshi]
command =  /usr/local/bin/uwsgi  /mnt/newdata/YunShi_API_Development/deploy/Yunshi_uwsgi.ini
directory = /mnt/newdata/root/YunShi_API_Development/
user = root
autostart = true
autorestart =  true
stdout_logfile = /mnt/newdata/log/Yunshi.log

stdout_logfile_maxbytes = 20MB  ; stdout 日志文件大小，默认 50MB
stdout_logfile_backups = 20     ; stdout 日志文件备份数 stdout 日志文件，需要注意当指定目录不存在时无法正常启动，所以需要手动创建目录（supervisord 会自动创建日志文件）
stderr_logfile = /mnt/newdata/log/Yunshi_access.log
loglevel=debug

[supervisord]
